# scripts/launchd/com.altaf.content-factory.open-content-prs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.altaf.content-factory.open-content-prs</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/altaf.shaikh/Code-Space/AI/content-factory/scripts/daily-content-sync.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/altaf.shaikh/Code-Space/AI/content-factory</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- 11:00 AM local time (Asia/Calcutta) daily -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>11</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/altaf.shaikh/Library/Logs/content-factory-open-content-prs.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/altaf.shaikh/Library/Logs/content-factory-open-content-prs.log</string>
</dict>
</plist>
